# etc/init.d/microsocks-pia
#!/sbin/openrc-run
description="MicroSocks SOCKS5 server bound to WireGuard (pia)"

depend() {
    need net
    after pia-port-forward
}

command="/usr/local/bin/microsocks"
command_args="-p 1080"
pidfile="/run/microsocks-pia.pid"
command_background="yes"
output_log="/var/log/microsocks.log"
error_log="/var/log/microsocks.err"

start_pre() {
    ebegin "Checking WireGuard (pia) status ..."
    ip addr show pia >/dev/null 2>&1
    if [ $? -ne 0 ]; then
        eend 1 "pia interface not up"
        return 1
    fi
    eend 0
}

start() {
    ebegin "Starting microsocks bound to WireGuard (pia)"
    start-stop-daemon --start --background \
        --make-pidfile --pidfile "$pidfile" \
        --exec "$command" -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping microsocks"
    start-stop-daemon --stop --pidfile "$pidfile"
    eend $?
}
