# etc/init.d/pia-port-forward
#!/sbin/openrc-run
# OpenRC service for PIA WireGuard + persistent port forwarding

description="PIA WireGuard VPN with persistent port forwarding"

depend() {
    need net
    after firewall
    before qbittorrent
}

command="/root/pia/run_pf_service.sh"
pidfile="/run/pia-port-forward.pid"
directory="/root/pia"

output_log="/var/log/pia-service.log"
error_log="/var/log/pia-service.err"

start_pre() {
    mkdir -p /var/log
    chmod 755 /var/log

    ebegin "Bringing up WireGuard interface: pia"
    wg-quick up pia >/var/log/wg-quick.log 2>&1
    eend $?
}

start() {
    ebegin "Starting PIA Port Forwarding Service"
    start-stop-daemon --start --background \
        --make-pidfile --pidfile "$pidfile" \
        --chdir "$directory" \
        --exec "$command" \
        --stdout "$output_log" --stderr "$error_log"
    eend $?
}

stop_pre() {
    ebegin "Stopping WireGuard interface: pia"
    wg-quick down pia >>/var/log/wg-quick.log 2>&1
    eend $?
}

stop() {
    ebegin "Stopping PIA Port Forwarding Service"
    start-stop-daemon --stop --pidfile "$pidfile"
    eend $?
}

restart() {
    ebegin "Restarting PIA Port Forwarding Service"
    svc_stop
    sleep 3
    svc_start
    eend $?
}
