# etc/init.d/pia-pf
#!/sbin/openrc-run

name="pia-pf"
description="PIA Port Forwarding screen session"

envdir="/opt/pia"
credfile="${envdir}/credentials.env"
optfile="${envdir}/options.env"

SCREEN="/usr/bin/screen"
SESSION="pia-pf"
LOGFILE="/var/log/pia-pf.log"
PIDFILE="/run/pia-pf.pid"

# Dependencies
depend() {
    need net
    need wg-pia
}

# Safely load KEY=VALUE lines from env files
load_env_file() {
    local file="$1"
    [ -f "$file" ] || return 0

    while IFS='=' read -r key value; do
        # Skip blank lines and comments
        [ -z "$key" ] && continue
        case "$key" in
            \#*) continue ;;
        esac
        # Skip keys with spaces
        echo "$key" | grep -q '[[:space:]]' && continue

        export "$key=$value"
    done < "$file"
}

start() {
    ebegin "Starting pia-pf (screen)"

    # Reset the log safely (truncate, don't delete)
    : > "$LOGFILE"

    # Load environment files
    load_env_file "$credfile"
    load_env_file "$optfile"

    # Ensure the script runs inside /opt/pia
    cd /opt/pia || exit 1

    # Kill stale screen session if it exists
    $SCREEN -S "$SESSION" -X quit 2>/dev/null || true

    # Start screen detached with stdin closed, stdout/stderr logged
    $SCREEN -L -Logfile "$LOGFILE" -dmS "$SESSION" sh -c "/opt/pia/run_vpn.sh </dev/null"

    # Record screen session PID for OpenRC
    sleep 1
    pid="$($SCREEN -ls | awk '/\.'"$SESSION"'\s/{print $1}' | cut -d. -f1)"
    [ -n "$pid" ] && echo "$pid" > "$PIDFILE"

    eend 0
}

stop() {
    ebegin "Stopping pia-pf"

    $SCREEN -S "$SESSION" -X quit 2>/dev/null || true
    rm -f "$PIDFILE"

    eend 0
}