# etc/init.d/pia-port
#!/sbin/openrc-run
#
# /etc/init.d/pia-port
# Extract PIA forwarded port and update qBittorrent config once on startup

name="pia-port"
description="Extract PIA forwarded port and push to qBittorrent once on startup"

pidfile="/run/pia-port.pid"
LOG="/var/log/pia-pf.log"
PORT_FILE="/run/pia-port.txt"
QB_PF_SCRIPT="/opt/pia/qbittorrent-pf.sh"

depend() {
    need pia-pf
    after pia-pf
}

start() {
    ebegin "Extracting PIA forwarded port"

    # Give PIA daemon time to initialize
    sleep 20

    if [ ! -f "$LOG" ]; then
        echo "PIA log not found: $LOG"
        eend 1
        return
    fi

    # Extract the port 7 lines from the end, strip ANSI sequences, keep only digits
    port=$(tail -n 7 "$LOG" | head -n 1 | sed -r 's/\x1B\[[0-9;]*[mK]//g' | tr -cd '0-9')

    if [ -z "$port" ]; then
        echo "Failed to extract port from log"
        eend 1
        return
    fi

    # Save to /run/pia-port.txt
    echo "$port" > "$PORT_FILE"
    echo "PIA forwarded port saved to $PORT_FILE: $port"

    # Run qbittorrent-pf script to update qBittorrent config
    if [ -x "$QB_PF_SCRIPT" ]; then
        "$QB_PF_SCRIPT"
        eend $?
    else
        echo "qBittorrent PF script not found or not executable: $QB_PF_SCRIPT"
        eend 1
    fi
}

stop() {
    ebegin "Stopping pia-port service"
    # Nothing to stop; one-shot
    eend 0
}