# etc/init.d/qbittorrent
#!/sbin/openrc-run

name="qbittorrent-nox"
description="qBittorrent-nox daemon"

command="/usr/bin/qbittorrent-nox"
command_args="--webui-port=8080"

pidfile="/run/qbittorrent.pid"
WEBUI="http://127.0.0.1:8080"
WAIT_SECS=30

depend() {
    need net
    use logger
    need wg-pia
    need pia-pf
}

start_pre() {
    mkdir -p /run
    chown qbittorrent:qbittorrent /run

    mkdir -p /home/qbittorrent/.config
    mkdir -p /home/qbittorrent/.cache
    mkdir -p /home/qbittorrent/.local/share
    mkdir -p /home/qbittorrent/Downloads

    chown -R qbittorrent:qbittorrent /home/qbittorrent
}

start() {
    ebegin "Starting $name"
    start-stop-daemon --start --background \
        --make-pidfile --pidfile $pidfile \
        --user qbittorrent \
        --env HOME="/home/qbittorrent" \
        --exec $command -- $command_args
    eend $?
}

start_post() {
    ebegin "Updating qBittorrent search plugins"

    for i in $(seq 1 $WAIT_SECS); do
        if curl -fs "$WEBUI/api/v2/app/version" >/dev/null; then
            curl -fs -X POST "$WEBUI/api/v2/search/updatePlugins" \
                && eend 0 && return
            break
        fi
        sleep 1
    done

    ewarn "WebUI not ready — plugin update skipped"
    eend 1
}

stop() {
    ebegin "Stopping $name"
    start-stop-daemon --stop --pidfile $pidfile
    eend $?
}