# etc/init.d/qbittorrent-clean-seeding
#!/sbin/openrc-run

name="qbt-remove-tv-seeding"
description="Remove seeding/uploading torrents in category tv"

pidfile="/run/qbt-remove-tv-seeding.pid"
command="/bin/sh"

depend() {
  need net
  need qbittorrent
  after qbittorrent
}

start() {
  ebegin "Starting qbt-remove-tv-seeding"

  start-stop-daemon \
    --start \
    --background \
    --make-pidfile \
    --pidfile "$pidfile" \
    --exec "$command" -- -c "
      exec >> /var/log/qbt-remove-tv-seeding.log 2>&1

      QBIT='http://localhost:8080'
      CATEGORY='tv'
      INTERVAL=300

      while true; do
        echo \"[\$(date)] Checking seeding torrents in category tv\"

        HASHES=\$(curl -s \"\$QBIT/api/v2/torrents/info?category=\$CATEGORY\" \
          | jq -r '
              [.[]
                | select(.state == \"uploading\" or .state == \"stalledUP\" or .state == \"forcedUP\")
                | .hash] | join(\"|\")
            ')

        if [ -n \"\$HASHES\" ]; then
          echo \"[\$(date)] Removing torrents: \$HASHES\"

          curl -s -X POST \
            --data-urlencode \"hashes=\$HASHES\" \
            --data \"deleteFiles=false\" \
            \"\$QBIT/api/v2/torrents/delete\"
        else
          echo \"[\$(date)] No seeding torrents found\"
        fi

        sleep \"\$INTERVAL\"
      done
    "

  eend $?
}

stop() {
  ebegin "Stopping qbt-remove-tv-seeding"
  start-stop-daemon --stop --pidfile "$pidfile"
  eend $?
}