# etc/init.d/qbittorrent-set-port
#!/sbin/openrc-run

description="Push PIA forwarded port to qBittorrent WebUI"

depend() {
    need net
    need qbittorrent
    after qbittorrent
}

pidfile="/run/qbittorrent-set-port.pid"

QBT_URL="http://localhost:8080"
LOG="/var/log/pia-pf.log"
PORT_FILE="/run/pia-port.txt"
INTERVAL=600
START_DELAY=15

start() {
    ebegin "Starting qbittorrent-set-port"

    (
        sleep "$START_DELAY"

        # Wait until qBittorrent API is ACTUALLY ready
        until curl -sf "$QBT_URL/api/v2/app/version" >/dev/null 2>&1; do
            sleep 2
        done

        last_port=""

        while true; do
            port=$(grep "Forwarded port" "$LOG" | tail -n1 \
                | sed -r 's/\x1B\[[0-9;]*[mK]//g' \
                | tr -d '\r' \
                | tr -cd '0-9')

            if [ -n "$port" ] && [ "$port" != "$last_port" ]; then
                echo "$port" > "$PORT_FILE"

                curl -i -X POST \
                    -d "json={\"listen_port\": $port}" \
                    "$QBT_URL/api/v2/app/setPreferences" \
                    >> /var/log/qbittorrent-set-port.log 2>&1

                last_port="$port"
            fi

            sleep "$INTERVAL"
        done
    ) &

    echo $! > "$pidfile"
    eend 0
}

stop() {
    ebegin "Stopping qbittorrent-set-port"

    if [ -f "$pidfile" ]; then
        kill "$(cat "$pidfile")" 2>/dev/null
        rm -f "$pidfile"
    fi

    eend 0
}
