# etc/init.d/qbittorrent-clean-seeding
#!/sbin/openrc-run

description="Remove seeding torrents in specific categories"

depend() {
    need net
    need qbittorrent
    after qbittorrent
}

pidfile="/run/qb-clean-seeding.pid"
LOG="/var/log/qb-clean-seeding.log"

QB_URL="http://localhost:8080"
CATEGORIES="movies tv games apps os music"
INTERVAL=600        # 10 minutes
STARTUP_DELAY=20

start() {
    ebegin "Starting qb-clean-seeding"

    (
        sleep "$STARTUP_DELAY"
        last_seen=""

        while true; do
            # Get all torrents
            torrents=$(curl -s "$QB_URL/api/v2/torrents/info")

            # Loop through categories
            for cat in $CATEGORIES; do
                hashes=$(echo "$torrents" \
                    | grep -Po '"hash":"\K[0-9a-f]+' \
                    | while read h; do
                        status=$(echo "$torrents" | grep -Po '"hash":"'"$h"'"(.*?)"state":"\K[^"]+')
                        category=$(echo "$torrents" | grep -Po '"hash":"'"$h"'"(.*?)"category":"\K[^"]*')
                        if [ "$category" = "$cat" ] && [ "$status" = "uploading" ]; then
                            echo "$h"
                        fi
                    done)

                for h in $hashes; do
                    curl -s -X POST \
                        -d "hashes=$h&deleteFiles=false" \
                        "$QB_URL/api/v2/torrents/delete" >> "$LOG" 2>&1
                    echo "$(date) - Removed seeding torrent $h in category $cat" >> "$LOG"
                done
            done

            sleep "$INTERVAL"
        done
    ) &

    echo $! > "$pidfile"
    eend 0
}

stop() {
    ebegin "Stopping qb-clean-seeding"

    if [ -f "$pidfile" ]; then
        kill "$(cat "$pidfile")" 2>/dev/null
        rm -f "$pidfile"
    fi

    eend 0
}