# etc/init.d/qbittorrent
#!/sbin/openrc-run

name="qbittorrent-nox"
description="qBittorrent-nox daemon"

command="/usr/bin/qbittorrent-nox"
command_args="--webui-port=8080"

pidfile="/run/qbittorrent.pid"

# Ensure correct service dependencies
depend() {
    need net
    use logger
    need wg-pia
    need pia-port
}

start_pre() {
    # Ensure runtime directory exists
    mkdir -p /run
    chown qbittorrent:qbittorrent /run

    # Ensure HOME and required dirs exist
    mkdir -p /home/qbittorrent/.config
    mkdir -p /home/qbittorrent/.cache
    mkdir -p /home/qbittorrent/.local/share
    mkdir -p /home/qbittorrent/Downloads

    chown -R qbittorrent:qbittorrent /home/qbittorrent
}

start() {
    ebegin "Starting $name"
    start-stop-daemon --start --background \
        --make-pidfile --pidfile $pidfile \
        --user qbittorrent \
        --env HOME="/home/qbittorrent" \
        --exec $command -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping $name"
    start-stop-daemon --stop --pidfile $pidfile
    eend $?
}