# etc/init.d/qbittorrent-pf
#!/sbin/openrc-run
# /etc/init.d/qbittorrent-pf
description="Start qBittorrent and set incoming port from PIA forwarded port"

command="/usr/bin/qbittorrent-nox"
pidfile="/run/qbittorrent-nox.pid"
output_log="/var/log/qbittorrent-pf.log"
error_log="/var/log/qbittorrent-pf.err"

depend() {
    need pia-port-forward
    after net pia-port-forward
}

get_port() {
    [ -f /tmp/forwarded_port ] || return 1
    PORT=$(cat /tmp/forwarded_port 2>/dev/null)
    [ -n "$PORT" ] || return 1
    echo "$PORT"
}

start() {
    ebegin "Waiting for valid PIA forwarded port"

    # Wait for /tmp/forwarded_port to exist AND be recent (within last 2 minutes)
    for i in $(seq 1 60); do
        if [ -f /tmp/forwarded_port ]; then
            age=$(( $(date +%s) - $(stat -c %Y /tmp/forwarded_port) ))
            if [ "$age" -lt 120 ]; then
                PORT=$(get_port)
                if [ -n "$PORT" ]; then
                    einfo "Got recent forwarded port: $PORT"
                    break
                fi
            fi
        fi
        sleep 2
    done

    if [ -z "$PORT" ]; then
        eerror "No recent forwarded port found!"
        return 1
    fi

    # Start qBittorrent if not running
    if pgrep -x qbittorrent-nox >/dev/null 2>&1; then
        einfo "qBittorrent already running"
    else
        ebegin "Starting qBittorrent (with umask 000)"
        start-stop-daemon --start --background \
            --make-pidfile --pidfile "$pidfile" \
            --exec /bin/sh -- -c "umask 000; exec $command" \
            --stdout "$output_log" --stderr "$error_log"
        eend $?
    fi

    # Wait for WebUI readiness
    einfo "Waiting for qBittorrent WebUI to respond..."
    for i in $(seq 1 60); do
        if curl -s http://127.0.0.1:8080/api/v2/app/version >/dev/null 2>&1; then
            break
        fi
        sleep 1
    done

    if ! curl -s http://127.0.0.1:8080/api/v2/app/version >/dev/null 2>&1; then
        eerror "qBittorrent WebUI not responding after 60s"
        return 1
    fi

    # Update port in qBittorrent
    einfo "Setting qBittorrent incoming port to $PORT"
    for attempt in $(seq 1 5); do
        if curl -s -X POST -d "json={\"listen_port\":$PORT}" \
            http://127.0.0.1:8080/api/v2/app/setPreferences >/dev/null 2>&1; then
            einfo "Successfully set qBittorrent port to $PORT"
            eend 0
            return 0
        fi
        sleep 3
    done

    eerror "Failed to set qBittorrent incoming port after 5 attempts"
    return 1
}

stop() {
    ebegin "Stopping qBittorrent"
    start-stop-daemon --stop --pidfile "$pidfile"
    eend $?
}
