# etc/init.d/qbt-set-books
#!/sbin/openrc-run

description="Set qBittorrent category to 'books' for torrents using myanonymouse.net"

depend() {
    need net
    need qbittorrent
    after qbittorrent
}

pidfile="/run/qb-set-books.pid"
LOG="/var/log/qb-set-books.log"

QB_URL="http://localhost:8080"
TARGET_TRACKER="myanonymouse.net"
CATEGORY_NAME="books"
INTERVAL=300
STARTUP_DELAY=20

start() {
    ebegin "Starting qb-set-books"

    (
        sleep "$STARTUP_DELAY"
        echo "$(date) - Service started" >> "$LOG"

        while true; do
            # Get all torrents
            torrents=$(curl -s "$QB_URL/api/v2/torrents/info")

            # Parse hashes of torrents with target tracker
            hashes=$(echo "$torrents" \
                | grep -Po '"hash":"\K[0-9a-f]+' \
                | while read h; do
                    tracker=$(echo "$torrents" | grep -Po '"hash":"'"$h"'"(.*?)"tracker":"\K[^"]+')
                    if echo "$tracker" | grep -q "$TARGET_TRACKER"; then
                        echo "$h"
                    fi
                done)

            # Set category for each matching torrent
            for h in $hashes; do
                curl -s -X POST \
                    -d "hash=$h&category=$CATEGORY_NAME" \
                    "$QB_URL/api/v2/torrents/setCategory" >> "$LOG" 2>&1
                echo "$(date) - Set torrent $h category to $CATEGORY_NAME" >> "$LOG"
            done

            sleep "$INTERVAL"
        done
    ) &

    echo $! > "$pidfile"
    eend 0
}

stop() {
    ebegin "Stopping qb-set-books"

    if [ -f "$pidfile" ]; then
        kill "$(cat "$pidfile")" 2>/dev/null
        rm -f "$pidfile"
    fi

    eend 0
}