# etc/init.d/qbt-set-books
#!/sbin/openrc-run

name="qbt-set-books"
description="Auto-assign MyAnonamouse torrents to books category"

pidfile="/run/qbt-set-books.pid"
command="/bin/sh"

depend() {
  need net
  need qbittorrent
  after qbittorrent
}

start() {
  ebegin "Starting qbt-set-books"

  start-stop-daemon \
    --start \
    --background \
    --make-pidfile \
    --pidfile "$pidfile" \
    --exec "$command" -- -c '
      exec >> /var/log/qbt-set-books.log 2>&1

      QBIT="http://localhost:8080"
      USER="admin"
      PASS="luigina"
      TRACKER="t.myanonamouse.net"
      CATEGORY="books"
      INTERVAL=300
      COOKIE="/run/qbt-set-books.cookies"

      while true; do
        echo "[$(date)] Checking torrents"

        curl -s -c "$COOKIE" -X POST \
          --data "username=$USER&password=$PASS" \
          "$QBIT/api/v2/auth/login" > /dev/null

        HASHES=""

        for HASH in $(curl -s -b "$COOKIE" "$QBIT/api/v2/torrents/info" | jq -r ".[].hash"); do
          if curl -s -b "$COOKIE" \
            "$QBIT/api/v2/torrents/trackers?hash=$HASH" \
            | jq -e --arg t "$TRACKER" ".[] | select(.url | contains(\$t))" > /dev/null
          then
            HASHES="$HASHES|$HASH"
          fi
        done

        HASHES="${HASHES#|}"

        if [ -n "$HASHES" ]; then
          curl -s -b "$COOKIE" -X POST \
            --data-urlencode "hashes=$HASHES" \
            --data-urlencode "category=$CATEGORY" \
            "$QBIT/api/v2/torrents/setCategory"
          echo "[$(date)] Updated category"
        else
          echo "[$(date)] No matching torrents"
        fi

        sleep "$INTERVAL"
      done
    '

  eend $?
}

stop() {
  ebegin "Stopping qbt-set-books"
  start-stop-daemon --stop --pidfile "$pidfile"
  eend $?
}