# etc/init.d/qbt-api-engine
#!/sbin/openrc-run

name="qbt-api-engine"
description="Config-driven qBittorrent API automation engine"

command="/usr/local/bin/qbt-api-engine"
pidfile="/run/qbt-api-engine.pid"

depend() {
  need net
  need qbittorrent
  after qbittorrent
}

start() {
  ebegin "Starting qbt-api-engine"
  start-stop-daemon \
    --start \
    --background \
    --make-pidfile \
    --pidfile "$pidfile" \
    --exec "$command"
  eend $?
}

stop() {
  ebegin "Stopping qbt-api-engine"
  start-stop-daemon \
    --stop \
    --pidfile "$pidfile"
  eend $?
}
